# deploy/com.citehustle.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.citehustle.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/__USER__/Github/cite-hustle/deploy/run_pipeline.sh</string>
        <string>incremental</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/__USER__/Library/Logs/cite-hustle/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/Library/Logs/cite-hustle/weekly.err.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
